# buildfiles/macosx/Library/LaunchDaemons/io.appserver.appserver-php5-fpm.plist
<?xml version='1.0' encoding='UTF-8'?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd" >
<plist version='1.0'>
    <dict>
        <key>Label</key>
        <string>io.appserver.appserver-php5-fpm</string>
        <key>ProgramArguments</key>
        <array>
            <string>/opt/appserver/sbin/php-fpm</string>
            <string>--pid=/opt/appserver/var/run/appserver-php5-fpm.pid</string>
            <string>--nodaemonize</string>
        </array>
        <key>Debug</key>
        <false/>
        <key>RunAtLoad</key>
        <true/>
        <key>OnDemand</key>
        <false/>
        <key>KeepAlive</key>
        <false/>
	</dict>
</plist>