# buildfiles/fedora/lib/systemd/system/watcher.service

[Unit]
Description=appserver watcher

[Service]
ExecStart=/opt/appserver/bin/watcher -D FOREGROUND

[Install]
WantedBy=multi-user.target
